# pub/batchacct-pub.cron
# Daily publish accounting records to APEL

42 9 * * * root python /usr/lib/python2.4/site-packages/batchacct/join.py --acctdbfile /etc/batchacct/connection --conf /etc/batchacct/pub --vofile /etc/batchacct/vos --ssm /opt/apel/ssm/messages/outgoing --logfile /var/log/batchacct/batchacct-pub.log
