# loccol/batchacct-partition.cron
# Monthly check DB table partitioning and add extra partitions if needed

42 10 17 * * root python /usr/lib/python2.4/site-packages/batchacct/partition.py --connfile /etc/batchacct/connection --table loc --logfile /var/log/batchacct/batchacct-partition.log
42 11 17 * * root python /usr/lib/python2.4/site-packages/batchacct/partition.py --connfile /etc/batchacct/connection --table ce --logfile /var/log/batchacct/batchacct-partition.log
